# salt/bccvl/plone_supervisor.ini
{% set plone_bin = '/home/plone/bccvl_buildout/bin' %}

[program:zeo]
command = {{ plone_bin }}/zeoserver fg
process_name = zeo
directory = {{ plone_bin }}
priority = 10
redirect_stderr = true
user = plone

{% for service in pillar['plone']['instances'].keys() %}
[program:{{ service }}]
command = {{ plone_bin }}/{{ service }} console
process_name = {{ service }}
directory = {{ plone_bin }}
priority = 20
redirect_stderr = true
user = plone
{% if pillar['plone']['instances'][service].get('debug', False) %}
autostart = false
{% endif %}

{% endfor %}

{% for service in pillar['plone']['workers'].keys() %}
[program:{{ service }}]
command = {{ plone_bin }}/{{ service }} console
process_name = {{ service }}
directory = {{ plone_bin }}
priority = 20
redirect_stderr = true
user = plone
{% if pillar['plone']['workers'][service].get('debug', False) %}
autostart = false
{% endif %}

{% endfor %}

# TODO: example crash listener should go into main supervisor config
# TODO: general supervisord setup like control interfaces and access credentials
# [eventlistener:crashmail]
# command = /home/plone/bccvl_buildout/bin/crashmail -p cache -m "root@localhost"
# events = PROCESS_STATE
# process_name=crashmail
# environment=SUPERVISOR_USERNAME='admin',SUPERVISOR_PASSWORD='neumideomruawbne',SUPERVISOR_SERVER_URL='http://127.0.0.1:9001'